# deb-src/etc/cron.d/scop-rest-api
* * * * * root run-one /opt/scop/bin/scop-rest-api >> /var/log/scop-rest-api.log 2>&1
